# debian/vyatta-rest.postinst
#!/bin/bash

prefix=/opt/vyatta

mkdir -p $prefix/config/tmp
chmod 2775 $prefix/config/tmp

# enable CGI for lighttpd
ln -sf /etc/lighttpd/conf-available/10-cgi.conf /etc/lighttpd/conf-enabled/

#vtysh access
rm -f /usr/lib/cgi-bin/vtysh

# enable https
ln -sf /etc/lighttpd/conf-available/10-ssl.conf /etc/lighttpd/conf-enabled/

# disable listing of web directories
if ! grep -q '^server.dir-listing' /etc/lighttpd/lighttpd.conf; then
  echo 'server.dir-listing = "disable"' \
    >>/etc/lighttpd/lighttpd.conf
else 
  sed -i 's/^server.dir-listing = .*/server.dir-listing = "disable"/' \
    /etc/lighttpd/lighttpd.conf
fi

# use our old document root for now
sed -i 's/^\(server\.document-root  *=\).*/\1 "\/var\/www"/' \
    /etc/lighttpd/lighttpd.conf

# disable listing of cgi-bin directory
chmod 711 /usr/lib/cgi-bin


# fix debian cgi config file
sed -i '/^\$HTTP.*127\.0\.0\.1/,/^}$/{s/^\(.*\)$/# \1/}' \
  /etc/lighttpd/conf-enabled/10-cgi.conf
if ! grep -q '^alias.url' /etc/lighttpd/conf-enabled/10-cgi.conf; then 
  echo 'alias.url += ( "/cgi-bin/" => "/usr/lib/cgi-bin/" )' \
    >>/etc/lighttpd/conf-enabled/10-cgi.conf
fi

# enable mod_redirect
if grep -q '^#.*mod_redirect' /etc/lighttpd/lighttpd.conf; then
  sed -i 's/^#\(.*mod_redirect.*\)$/ \1/' /etc/lighttpd/lighttpd.conf
fi

# enable ipv6 support
if ! grep -q '\[::\]' /etc/lighttpd/lighttpd.conf; then
  {
    echo -e '\n# Listen to IPv6 HTTP'
    echo '$SERVER["socket"] == "[::]:80" {  }'
    echo '# Listen to IPv6 HTTPS'
    echo '$SERVER["socket"] == "[::]:443" { include "conf-enabled/10-ssl.conf"}'
  } >> /etc/lighttpd/lighttpd.conf
fi

# enable syslog for lighttpd
if ! grep -q '^server.errorlog-use-syslog' /etc/lighttpd/lighttpd.conf; then
  {
    echo -e '\n#/usr/pkg/etc/lighttpd/lighttpd.conf'
    echo 'server.modules += ( "mod_accesslog" )'
    echo 'server.errorlog-use-syslog = "enable"'
    echo 'accesslog.use-syslog = "enable"'
  } >> /etc/lighttpd/lighttpd.conf
fi

# enable syslog for lighttpd
if ! grep -q '^server.event-handler =' /etc/lighttpd/lighttpd.conf; then
  echo 'server.event-handler = "linux-sysepoll"' >> /etc/lighttpd/lighttpd.conf
fi

#remove tilde access-deny
if grep -q '^url.access-deny' /etc/lighttpd/lighttpd.conf; then
  sed -i 's/^\(.*url.access-deny.*\)/# \1/' /etc/lighttpd/lighttpd.conf
fi

#remove default ipv6 support with fallback handling to ipv4
if grep -q 'use-ipv6.pl' /etc/lighttpd/lighttpd.conf; then
  sed -i 's/^\(.*use-ipv6.pl.*\)/# \1/' /etc/lighttpd/lighttpd.conf
fi

#listen on 443 by default
if grep -q 'use-ipv6.pl' /etc/lighttpd/lighttpd.conf; then
  sed -i 's/^\(server\.port.*\)= 80/\1 = 443/' /etc/lighttpd/lighttpd.conf
fi


# redirect http to https
if ! grep -q '^# redirect all http to https$' /etc/lighttpd/lighttpd.conf; then
  cat <<'EOF' >>/etc/lighttpd/lighttpd.conf

# redirect all http to https
$HTTP["scheme"] == "http" {
  $HTTP["host"] =~ "^(.*)$" {
    url.redirect = (
      "^(/Vyatta/.*)$" => "https://%1$1",
      "^(/ext-2.2/.*)$" => "https://%1$1"
    )
  }
}
EOF
fi

# XXX www-data user needs to be in these groups for now
usermod -a -G "adm,sudo,users,routeadm,vyattacfg" www-data

# remove init links
update-rc.d -f lighttpd disable
update-rc.d -f lighttpd remove

#enable fastcgi
  cat <<'EOF' >/etc/lighttpd/conf-enabled/10-fastcgi.conf
server.modules   += ( "mod_fastcgi" )

## Start an FastCGI server for php (needs the php5-cgi package)
$HTTP["url"] =~ "/rest/" {
fastcgi.server    = ( "" => 
((
    "socket" => "/tmp/fcgi.socket",
    "max-procs" => 2,
    "idle-timeout" => 20,
    "bin-environment" => ( 
    "PHP_FCGI_CHILDREN" => "4",
    "PHP_FCGI_MAX_REQUESTS" => "10000"
  ),
    "bin-copy-environment" => (
    "PATH", "SHELL", "USER"
  ),
  "broken-scriptfilename" => "enable"
))
)
}
EOF


# fix debian cgi config file
if ! grep -q '(/Vyatta2/.*)$" => "https://%1$1",$' /etc/lighttpd/lighttpd.conf; then
  sed  -i 's/      "^(\/Vyatta\/.*)$" => "https:\/\/%1$1",/      "(\/Vyatta\/.*)$" => "https:\/\/%1$1",\n      "(\/Vyatta2\/.*)$" => "https:\/\/%1$1",/'   /etc/lighttpd/lighttpd.conf
fi

chmod a+w /opt/vyatta/tmp/webgui2

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

